# install_apache_ignite/templates/ignite-server.service
[Unit]
Description=Apache Ignite server
After=network.target
Requires=network.target

[Service]

